# totp-autobackup/totp-autobackup.timer
[Unit]
Description=Totp-autobackup timer

[Timer]
AccuracySec=10s
Unit=totp-autobackup.service

[Install]
WantedBy=default.target
